# packaging/debian/postinst
#!/bin/bash
set -e

APP_DIR="/opt/leiturgia"
APP_USER="pi"

cd "${APP_DIR}"

python3 -m venv .venv
.venv/bin/pip install --upgrade pip -q
.venv/bin/pip install -r requirements.txt -q

mkdir -p \
  media/images \
  media/videos \
  data/lyrics \
  output

if [ ! -f data/program.json ]; then
  cat > data/program.json <<'EOF'
{
  "church": "",
  "date": "",
  "pianist": "",
  "song_leader": "",
  "service_programs": [
    {
      "id": "sample-program",
      "name": "Sample Program",
      "time": "",
      "items": [
        {"item_id": "sp-001", "type": "participant", "title": "Opening Prayer", "part": "Opening Prayer", "participant": ""},
        {"item_id": "sp-002", "type": "song",        "title": "Opening Song",   "hymn_number": ""},
        {"item_id": "sp-003", "type": "media",       "title": "Welcome",        "media_type": "image",    "url": "/media/images/leiturgia-welcome.png"},
        {"item_id": "sp-004", "type": "media",       "title": "Sample Video",   "media_type": "video",    "url": "", "autoplay": true, "loop": false, "mute": false},
        {"item_id": "sp-005", "type": "content",     "title": "Announcements",  "content": ""}
      ]
    }
  ],
  "service_team": []
}
EOF
  chown "${APP_USER}:${APP_USER}" data/program.json
fi

if [ ! -f config.json ]; then
  SECRET=$(openssl rand -hex 32)
  cat > config.json <<EOF
{
  "pin": "1234",
  "session_secret": "${SECRET}",
  "session_timeout_hours": 8,
  "max_login_attempts": 5
}
EOF
  chown "${APP_USER}:${APP_USER}" config.json
  chmod 640 config.json
fi

chown -R "${APP_USER}:${APP_USER}" "${APP_DIR}"

DESKTOP_DIR="/home/${APP_USER}/Desktop"
SHORTCUT="${DESKTOP_DIR}/Leiturgia.desktop"
if [ -d "${DESKTOP_DIR}" ]; then
  cat > "${SHORTCUT}" <<EOF
[Desktop Entry]
Version=1.0
Type=Application
Name=Leiturgia
Comment=Open Leiturgia operator console
Exec=sh -c 'chromium-browser http://localhost:5000 || chromium http://localhost:5000 || xdg-open http://localhost:5000'
Icon=${APP_DIR}/static/leiturgia-icon.png
Terminal=false
Categories=Application;
EOF
  chmod +x "${SHORTCUT}"
  chown "${APP_USER}:${APP_USER}" "${SHORTCUT}"
fi

systemctl daemon-reload
systemctl enable leiturgia
systemctl start leiturgia
