# packaging/debian/prerm
#!/bin/bash
set -e

systemctl stop leiturgia leiturgia-display 2>/dev/null || true
systemctl disable leiturgia leiturgia-display 2>/dev/null || true
